# deploy/log-guardian.service.d/20-readwrite.conf
# Webhook metrics + state — ProtectSystem=strict altinda yazma izni
[Service]
ReadWritePaths=/etc/log-guardian /var/log/nginx /var/lib/log-guardian
