# deploy/log-guardian.service.d/30-password-file.conf
# Servis parolasi — env/KDF uyumsuzlugunda service.password kullanilir
[Service]
ExecStart=
ExecStart=/usr/local/bin/log-guardian /var/log/nginx/access.log \
    --rules /etc/log-guardian/rules.conf \
    --db /etc/log-guardian/events.db \
    --password-file /etc/log-guardian/service.password \
    --follow --no-tui
